# scripts/dwm-panel.service
[Unit]
Description=DWM panel deamon

[Service]
Type=simple
#PIDFile=/var/run/dwm-panel.pid
ExecStart=/usr/local/bin/dwm-panel

